# .config/systemd/user/blueman-applet.service
[Unit]
Description=Blueman Applet


[Service]
ExecStart=/usr/bin/blueman-applet
Restart=on-failure

[Install]
WantedBy=graphical-session.target

