# .config/systemd/user/graphical-session.target
[Unit]
Description=Graphical session (Qtile)
